# etc/autostart/piday.desktop
#!/bin/bash

[Desktop Entry]
Exec= lxterminal -e "/home/pi/srp/piday/etc/start.sh" 
Type=Application;


